# deploy/stock-backup.timer
# stock-backup.timer — run the database backup every night at ~02:15, before the
# 02:30 reconcile check.
#
# Persistent=true means a run missed while the machine was off fires at the next
# boot, so a backup is never silently skipped. Enable with:
#   systemctl enable --now stock-backup.timer
# Inspect with:  systemctl list-timers stock-backup.timer

[Unit]
Description=Run Stock Manager database backup nightly

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
